# systemd/claude-relay.service
[Unit]
Description=Claude Code Relay Daemon
After=network.target

[Service]
Type=simple
ExecStart=%h/.npm-global/bin/claude-relay start
ExecStop=%h/.npm-global/bin/claude-relay stop
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=default.target
